# ansible/playbooks/templates/lynis-audit-monitor.service.j2
[Unit]
Description=Lynis System Security Audit Monitor
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/lynis-audit-monitor.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lynis-audit-monitor

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/log /var/lib/lynis-audit-monitor /var/lock /etc/security-monitor
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=no
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

[Install]
WantedBy=multi-user.target
