# ansible/playbooks/templates/ssh-security-monitor.service.j2
[Unit]
Description=SSH Security Monitor
After=network.target
Docs=https://github.com/forwardemail/forwardemail.net

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ssh-security-monitor.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ssh-security-monitor
User=root
Group=root

[Install]
WantedBy=multi-user.target
