# ansible/playbooks/templates/ssh-security-monitor.timer.j2
[Unit]
Description=SSH Security Monitor Timer
Documentation=https://github.com/forwardemail/forwardemail.net
Requires=ssh-security-monitor.service

[Timer]
# Run every N minutes (configurable via SSH_MONITOR_INTERVAL env var, default 10)
OnBootSec=3min
OnUnitActiveSec={{ lookup('env', 'SSH_MONITOR_INTERVAL') | default('10', true) }}min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
